# files/systemd-first-run.service.in
# SPDX-FileCopyrightText: 2025 Kristen McWilliam <kristen@kde.org>
#
# SPDX-License-Identifier: LGPL-2.1-or-later

# systemd service file for loading the KDE Initial System Setup application
#
# systemd/system/kde-initial-system-setup.service

[Unit]
Description=KDE Initial System Setup
After=systemd-user-sessions.service plymouth-quit.service systemd-logind.service
Before=display-manager.service

[Service]
Type=oneshot
User=kde-initial-system-setup
ExecStart=${CMAKE_INSTALL_FULL_LIBEXECDIR}/kde-initial-system-setup-bootutil

[Install]
WantedBy=multi-user.target
